# finedge.service
[Unit]
Description=FinEdge AI Edge Predictor
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/finedge
ExecStart=/opt/finedge/venv/bin/python3 app.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
